# macos/etc/inittab
::sysinit:/etc/init.d/rcS
hvc0::respawn:/bin/login root
::shutdown:/bin/umount -a -r
::restart:/sbin/init
